# runipfs.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit unlimited

exec sudo -u www-data HOME=/home/www-data ipfs daemon
